# terraform/modules/app/files/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=appuser
WorkingDirectory=/home/appuser/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always
